# led-matrix.service
[Unit]
Description=LED Matrix

[Service]
ExecStart=/home/pi/matrixpi/led_script.bash
Restart=on-failure

[Install]
WantedBy=multi-user.target
